# deploy/cpa-gateway.service
[Unit]
Description=CPA multi-key quota gateway (front of CLIProxyAPI)
After=network.target

[Service]
Type=simple
# Adjust paths to your install
WorkingDirectory=/opt/grok-regkit
Environment=CPA_GATEWAY_ROOT=/opt/cliproxyapi
Environment=CPA_GATEWAY_HOST=0.0.0.0
Environment=CPA_GATEWAY_PORT=8318
Environment=CPA_UPSTREAM=http://127.0.0.1:8317
# Shown when running: python cpa_gateway.py add
Environment=CPA_PUBLIC_BASE=https://api.example.com/cpa/v1
ExecStart=/usr/bin/python3 /opt/grok-regkit/deploy/cpa_gateway.py serve
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
